# ecagentd.service
[Unit]
Description=ECManaged Agent for monitoring and deployment

[Service]
Type=simple
PIDFile=/opt/ecmanaged/ecagent/twistd.pid
#ExecStart=/usr/bin/twistd -y /opt/ecmanaged/ecagent/ecagentd.tac
ExecStart=\
        /bin/twistd \
        --nodaemon \
        --pidfile=/opt/ecmanaged/ecagent/twistd.pid \
        --no_save \
        --python=/opt/ecmanaged/ecagent/ecagentd.tac

WorkingDirectory=/opt/ecmanaged/ecagent
Restart=always
RestartSec=300

[Install]
WantedBy=multi-user.target
