# build/debian/cron.d
@reboot /opt/ecamanged/ecagent/init check > /dev/null 2>&1
*/5 *	* * *	root	/opt/ecamanged/ecagent/init check > /dev/null 2>&1
