# cron.d/ecmanaged-ecagent
SHELL=/bin/bash
HOME=/opt/ecmanaged/ecagent/
*/2 *	* * *	root	/opt/ecmanaged/ecagent/init check > /dev/null 2>&1

