# cron.d/ecmanaged-ecagent
SHELL=/bin/bash
HOME=/opt/ecmanaged/ecagent/
*/1 *	* * *	root	/opt/ecmanaged/ecagent/init check
@reboot /opt/ecmanaged/ecagent/init check
