# cron.d/ecmanaged-ecagent
*/5 *	* * *	root	/opt/ecmanaged/ecagent/init check > /dev/null 2>&1
