# cron.d/ecmanaged-ecagent
SHELL=/bin/bash
HOME=/opt/ecmanaged/ecagent/
*/5 *	* * *	root	/opt/ecmanaged/ecagent/init check
