# build/redhat/etc/cron.d/ecmanaged-ecagent-systemd
*/5 * * * root  service ecagentd status || service ecagentd start
