# build/redhat/etc/cron.d/ecmanaged-ecagent-a
*/5 *	* * *	root	/etc/init.d/ecagentd check > /dev/null 2>&1
